# config/cctv-footage-api.service
[Unit]
Description=CCTV Dashboard Footage API
After=network-online.target
Wants=network-online.target
# Footage lives on the sshfs auto-mount at /media/share; make sure it is
# mounted before the API starts (triggers the x-systemd.automount unit).
RequiresMountsFor=/media/share

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/python3 __PROJECT_DIR__/scripts/footage_api.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
