# config/footage-api.service
[Unit]
Description=CCTV Dashboard Footage API
After=network.target

[Service]
Type=simple
User=__RUN_USER__
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/python3 __PROJECT_DIR__/scripts/footage_api.py
Restart=always
RestartSec=3

[Install]
WantedBy=multi-user.target
